# templates/grader/omegaup-grader.service.erb
[Unit]
Description=omegaUp grading service
<% if @local_database %>
Requires=mysql.service
After=network.target mysql.service
<% end %>

[Service]
Type=simple
User=omegaup
Group=omegaup
ExecStart=/usr/bin/omegaup-grader
WorkingDirectory=/var/lib/omegaup
Restart=always

[Install]
WantedBy=multi-user.target
